# injected-units/setup-lsi-cards-in-raid0.service
[Unit]
ConditionPathExists=!/etc/%m.raid-setup

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/bin/lsi-cards-make-raid0 %m